# systemd/pearl-miner.service
[Unit]
Description=Pearl Miner
After=network-online.target gost-client.service
Wants=network-online.target
Requires=gost-client.service

[Service]
Type=simple
EnvironmentFile=/etc/gost-thread/miner.env
ExecStartPre=/usr/local/lib/gost-thread/wait-for-pearl-miner-pool
ExecStart=/usr/local/lib/gost-thread/run-pearl-miner
Restart=always
RestartSec=5
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
