# systemd/pearl-miner.service
[Unit]
Description=Pearl Miner
After=network-online.target gost-client.service
Wants=network-online.target
Requires=gost-client.service

[Service]
Type=simple
EnvironmentFile=/etc/gost-thread/miner.env
ExecStartPre=/usr/local/lib/gost-thread/wait-for-pearl-miner-pool
ExecStart=/bin/sh -c 'cd "${MINER_WORKDIR}" && if [ -n "${MINER_LD_LIBRARY_PATH:-}" ]; then export LD_LIBRARY_PATH="${MINER_LD_LIBRARY_PATH}:${LD_LIBRARY_PATH:-}"; fi && exec "${MINER_BIN}" ${MINER_ARGS}'
Restart=always
RestartSec=5
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
