# systemd/pearl-miner.service
[Unit]
Description=Pearl Miner
Conflicts=akoya-miner.service
After=network-online.target gost-client.service
Wants=network-online.target
Requires=gost-client.service

[Service]
Type=simple
EnvironmentFile=/etc/gost-thread/miner.env
ExecStartPre=/usr/local/lib/gost-thread/wait-for-pearl-miner-pool
ExecStart=/bin/sh -c 'cd "${MINER_WORKDIR}" && exec "${MINER_BIN}" ${MINER_ARGS}'
Restart=always
RestartSec=5
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
